# crontab
0 0 */2 * * root . /tmp/.docker.env; /usr/local/sbin/renew-certificate >> /var/log/renew-certificate.log
