# crontab
0 0 */2 * * /usr/local/sbin/renew-certificate >> /var/log/renew-certificate.log