# crontab
0 0 */2 * * root env >> /etc/environment && /usr/local/sbin/renew-certificate >> /var/log/renew-certificate.log
