# crontab
0 0 */2 * * root . $HOME/.profile; /usr/local/sbin/renew-certificate >> /var/log/renew-certificate.log
