# crontab
0 0 */2 * * root env - `cat  /tmp/.docker.env` /usr/local/sbin/renew-certificate >> /var/log/renew-certificate.log
